# bin/portal@.service
[Unit]
Description=portal server
After=network.target

[Service]
Type=forking
PidFile=/home/ubuntu/portal-%i/tmp/pid

User=ubuntu
Group=ubuntu

ExecStart=/bin/bash -c /home/ubuntu/portal-%i/bin/swfstart
SuccessExitStatus=SIGKILL

StartLimitBust=2
StartLimitIntervalSec=300


TimeoutSec=300

Restart=always

PrivateTmp=true

[Install]
WantedBy=multi-user.target

